# deploy/systemd/wuwa-echo-backend.service
[Unit]
Description=Wuwa Echo Backend Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/wuwa/echo/backend
EnvironmentFile=-/root/wuwa/echo/backend/.env
ExecStart=/root/wuwa/echo/backend/.venv/bin/python -m uvicorn main:app --host 0.0.0.0 --port 8888
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
